# clawd/scripts/com.clawdbot.auto-resume.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clawdbot.auto-resume</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/felipevieira/.clawdbot/scripts/auto-resume.sh</string>
        <string>check</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/Users/felipevieira/.clawdbot/logs/auto-resume-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/felipevieira/.clawdbot/logs/auto-resume-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
